# supervisor/rq_worker_low.conf
; ~/projects/deepseek_dispatcher-new/supervisor/rq_worker_low.conf
[program:rq_worker_low]
command=/usr/local/bin/rq worker low --with-scheduler
directory=/app
autostart=true
autorestart=true
; 以下两行被注释掉，让 RQ Worker 的日志由 Python 的 logging 模块直接写入文件
; stderr_logfile=/app/logs/worker_low.err.log
; stdout_logfile=/app/logs/worker_low.out.log
environment=PYTHONPATH="/app"